# cookbooks/cephco-generic/files/default/ceph-libvirt-dns.conf
description "Libvirt/PowerDNS integration"

start on runlevel [2345]
stop on runlevel [!2345]

console log

respawn

exec /srv/inktank.com/ceph-libvirt-dns/ceph-libvirt-dns.py

